# deploy/systemd/pixelwise-deploy.service
[Unit]
Description=PixelWise auto-deploy
After=network-online.target

[Service]
Type=oneshot
User=produser
WorkingDirectory=/opt/pixelwise
ExecStart=/opt/pixelwise/deploy/auto-deploy.sh
